# playbooks/roles/ansible-role-rbenv/templates/etc/profile.d/rbenv.sh.j2
export PATH={{ rbenv_root_path }}/bin:$PATH
export RBENV_ROOT="{{ rbenv_root_path }}"
eval "$(rbenv init -)"